# ansible/roles/user_auth/templates/auth_api.service.j2
[Unit]
Description=User Authentication API
After=network.target

[Service]
User=api
WorkingDirectory=/path/to/your/api
ExecStart=/usr/bin/python3 /path/to/your/api/app.py
Restart=always

[Install]
WantedBy=multi-user.target
